# solvable/supervisor/eventhandler.conf
[program:eventhandler]
user=root
directory=/.eventhandler
command=node eventHandlers.js
