# solvable/supervisor/eventhandler.conf
[program:eventhandler]
user=root
directory=/.tutorial
command=node eventHandlers.js
