# solvable/supervisor/mongo.conf
[program:mongod]
command=/usr/local/bin/mongod
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
user=mongodb
startretries=10
priority=10
